# sota_client.service
[Unit]
Description=Sota Client
After=network.target
Requires=dbus.service
After=dbus.service

[Service]
Type=simple
ExecStart=/opt/ats/sota_client/bin/sota_client -c /opt/ats/sota_client/conf/client.toml

[Install]
WantedBy=multi-user.target
